# deploy/gridai.service
[Unit]
Description=GridAI Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/gridai
Environment=NODE_ENV=production
ExecStart=/usr/bin/env bash -lc "source /opt/nvm/nvm.sh && nvm use 22 && npm start"
Restart=always
RestartSec=5
StandardOutput=append:/var/log/gridai.log
StandardError=append:/var/log/gridai.log

[Install]
WantedBy=multi-user.target
